# service/subversion/run
#!/bin/sh

exec chpst -u svn svnserve --daemon --foreground -r /srv/svn
